# root/etc/systemd/system/dnscrypt-proxy.socket.d/override.conf
# https://github.com/jedisct1/dnscrypt-proxy/pull/416
[Socket]
ListenStream=
ListenDatagram=
# note that avahi uses this port
ListenStream=127.0.0.1:5353
ListenDatagram=127.0.0.1:5353
