# common/.config/systemd/user/udiskie.service

[Unit]
Description=Udiskie Automounting Daemon

[Service]
ExecStart=/usr/bin/env udiskie --quiet --tray

[Install]
WantedBy=default.target
