# root/etc/systemd/system/connman.service.d/override.conf
# https://wiki.archlinux.org/index.php/DNSCrypt
[Service]
ExecStart=
# don't alter resolv.conf or act as a DNS proxy
ExecStart=/usr/bin/env connmand -n --nodnsproxy
