# root/etc/systemd/system/connman.service.d/override.conf
# https://wiki.archlinux.org/index.php/DNSCrypt
[Service]
ExecStart=
ExecStart=/usr/bin/connmand -n --nodnsproxy
