# common/.config/systemd/user/emacs.service.d/override.conf
[Service]
# workaround for fcitx (and others) to work in emacs
Environment="LC_CTYPE=ja_JP.UTF-8"
